# enforce-mig-slices.cron
# Enforce MIG slices at every boot if MIG mode is enabled
@reboot  root  sleep 90 && /usr/local/bin/enforce-mig-slices.sh >> /var/log/enforce-mig-slices.log 2>&1
